# projects/drop/drop.service
[Unit]
Description=Drop sharing server
After=network.target

[Service]
Type=simple
User=epatel
WorkingDirectory=/home/epatel/vps-ai/projects/drop
ExecStart=/home/epatel/vps-ai/projects/drop/venv/bin/python server.py --port 8090
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
